# deploy/odaily-jin10-monitor.service
[Unit]
Description=OdAIly Jin10 monitor worker
After=network-online.target odaily-local-pipeline.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/OdAIly
EnvironmentFile=-/opt/OdAIly/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/OdAIly/.venv/bin/python backend/src/main.py jin10-monitor-worker
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
