# deploy/odaily-pipeline-supervisor.service
[Unit]
Description=OdAIly pipeline supervisor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/OdAIly
EnvironmentFile=-/opt/OdAIly/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/OdAIly/.venv/bin/python backend/src/main.py pipeline-supervisor
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
